# dist/bluesky-labeler-logger.target
[Unit]
Description=Bluesky Labeler Logger instances